# deploy/skgt-vote.service
[Unit]
Description=SKGT Vote - 실시간 투표 & 채팅 (Next.js + Socket.io)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/skgt-vote

# 정적 환경변수
Environment=NODE_ENV=production
Environment=PORT=3000
# 비밀 환경변수(ADMIN_SLUG 등)는 .env 에서 로드. CI가 GitHub Secret 으로 매 배포 시 갱신.
# '-' 접두사: 파일이 없어도 서비스는 기동(그 경우 ADMIN_SLUG 기본값 "dev").
EnvironmentFile=-/home/ubuntu/skgt-vote/.env

# tsx 로 커스텀 서버(server.ts) 직접 구동 (npm/cross-env 불필요)
ExecStart=/home/ubuntu/skgt-vote/node_modules/.bin/tsx server.ts

# 비정상 종료 시 자동 재시작
Restart=always
RestartSec=3

# 로그는 journald 로 (journalctl -u skgt-vote 로 확인)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skgt-vote

[Install]
WantedBy=multi-user.target
